# supervisord.conf
[supervisord]
nodaemon=true

[program:dvna]
priority=10
user=app
directory=/app/dvna
command=node dvna.js

[program:gruyere]
priority=10
user=app
directory=/app/gruyere
command=python gruyere.py
